# kimi_escrow.service
[Unit]
Description=KIMI Escrow Django Application
Documentation=https://github.com/your-org/kimi-escrow
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=notify
User=kimi_escrow
Group=kimi_escrow
WorkingDirectory=/var/www/kimi-escrow
Environment=PATH=/var/www/kimi-escrow/venv/bin
Environment=DJANGO_SETTINGS_MODULE=kimi_escrow.settings_production
Environment=PYTHONPATH=/var/www/kimi-escrow

# Configuration Gunicorn
ExecStart=/var/www/kimi-escrow/venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --worker-connections 1000 \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/www/kimi-escrow/logs/gunicorn_access.log \
    --error-logfile /var/www/kimi-escrow/logs/gunicorn_error.log \
    --log-level info \
    --timeout 120 \
    --keep-alive 2 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --preload \
    --user kimi_escrow \
    --group kimi_escrow \
    kimi_escrow.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/kimi-escrow/logs /var/www/kimi-escrow/media /var/www/kimi-escrow/backups
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kimi-escrow

[Install]
WantedBy=multi-user.target
